# scripts/module-boot.conf
description "ROS scripts module"

start on net-device-added INTERFACE=eth1
stop on shutdown #runlevel [!2345] 

respawn

exec /home/oddbot/catkin_ws/src/locomotion_module/scripts/module.bash
